# project/etc/supervisor.d/laravel_worker.ini
[program:laravel_worker]
command=php artisan queue:work --tries=3
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=3
priority = 900
redirect_stderr=true
stderr_logfile=%(here)s/../storage/logs/%(program_name)s.log
stdout_logfile=%(here)s/../storage/logs/%(program_name)s.log
